# update-motd.d/55-ups
#!/bin/bash
export TERM=xterm-256color
#printf "\n"
printf "$(tput setaf 86)" #main text color


printf "$(tput sgr0)" #reset color
